# templates/default/systemd.service.erb
[Unit]
Description=<%= @description %>
<% if @requires -%>
Requires=<%= @requires %>
<% end %>
<% if @after -%>
After=<%= @after %>
<% end %>
<% if @bindsto -%>
BindsTo=<%= @bindsto %>
<% end %>
<% if @partof -%>
PartOf=<%= @partof %>
<% end %>
<% if @wants -%>
Wants=<%= @wants %>
<% end %>
<% if @before -%>
Before=<%= @before %>
<% end %>

[Service]
<% if @user -%>
User=<%= @user %>
<% end %>
<% if @timeoutstartsec -%>
TimeoutStartSec=<%= @timeoutstartsec %>
<% end %>
<% if @execstartpre -%>
ExecStartPre=<%= @execstartpre %>
<% end %>
<% @environment.each do |key, value| -%>
Environment=<%= key %>=<%= value %>
<% end %>
<% if @execstartpost -%>
ExecStartPost=<%= @execstartpost %>
<% end %>
<% if @execstoppost -%>
ExecStopPost=<%= @execstoppost %>
<% end %>
<% if @execreload -%>
ExecReload=<%= @execreload %>
<% end %>
<% if @killmode -%>
KillMode=<%= @killmode %>
<% end %>
<% if @restart -%>
Restart=<%= @restart %>
<% end %>
<% if @execstop -%>
ExecStop=<%= @execstop %>
<% end %>
ExecStart=<%= @execstart %>

[Install]
WantedBy=multi-user.target
