# templates/default/systemd.service.erb
[Unit]
Description=<%= @description %>
<% @requires.each do |value| -%>
Requires=<%= value %>
<% end %>
<% @after.each do |value| -%>
After=<%= value %>
<% end %>
<% @bindsto.each do |value| -%>
BindsTo=<%= value %>
<% end %>
<% @partof.each do |value| -%>
PartOf=<%= value %>
<% end %>
<% @wants.each do |value| -%>
Wants=<%= value %>
<% end %>
<% if @before -%><% @before.each do |value| -%>
Before=<%= value %>
<% end %><% end %>

[Service]
<% if @user -%>
User=<%= @user %>
<% end %>
<% if @timeoutstartsec -%>
TimeoutStartSec=<%= @timeoutstartsec %>
<% end %>
<% @execstartpre.each do |value| -%>
ExecStartPre=<%= value %>
<% end %>
<% @environment.each do |key, value| -%>
Environment=<%= key %>=<%= value %>
<% end %>
<% @execstartpost.each do |value| -%>
ExecStartPost=<%= value %>
<% end %>
<% @execstoppost.each do |value| -%>
ExecStopPost=<%= value %>
<% end %>
<% if @execreload -%>
ExecReload=<%= @execreload %>
<% end %>
<% if @killmode -%>
KillMode=<%= @killmode %>
<% end %>
<% if @restart -%>
Restart=<%= @restart %>
<% end %>
<% if @execstop -%>
ExecStop=<%= @execstop %>
<% end %>
ExecStart=<%= @execstart %>

[Install]
WantedBy=multi-user.target

<% if @xfleet.length > 0 -%>
[X-Fleet]
<% end %>
<% @xfleet.each do |key, value| -%>
<%= key %>=<%= value %>
<% end %>
